# src/templates/bt-autoconnect.service
[Unit]
Description=Reconnexion automatique enceinte Bluetooth
After=bluetooth.service
Wants=bluetooth.service

[Service]
Type=simple
EnvironmentFile=${INSTALL_DIR}/soundspot.conf
ExecStartPre=/bin/sleep 10
ExecStart=${INSTALL_DIR}/bt-connect.sh
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
