# src/templates/soundspot-client-satellite.service
[Unit]
Description=SoundSpot Satellite — Snapclient
After=bluetooth.service network-online.target
Wants=bluetooth.service network-online.target

[Service]
Type=simple
EnvironmentFile=${INSTALL_DIR}/soundspot.conf
Environment="XDG_RUNTIME_DIR=/run/user/1000"
ExecStartPre=/bin/sleep 15
ExecStart=/usr/bin/snapclient -h ${MASTER_HOST} -p ${SNAPCAST_PORT} --player pulse:buffer_time=50
Restart=always
RestartSec=10
User=${SOUNDSPOT_USER}
StandardOutput=journal

[Install]
WantedBy=multi-user.target
