# src/templates/soundspot-client-master.service
[Unit]
Description=SoundSpot — Snapclient local (BT speaker)
After=snapserver.service bluetooth.service time-sync.target user@${SOUNDSPOT_UID}.service
Wants=snapserver.service bluetooth.service

[Service]
Type=simple
User=${SOUNDSPOT_USER}
Group=${SOUNDSPOT_USER}
Environment="XDG_RUNTIME_DIR=/run/user/${SOUNDSPOT_UID}"
Environment="PULSE_SERVER=unix:/run/user/${SOUNDSPOT_UID}/pulse/native"
ExecStartPre=${INSTALL_DIR}/wait-pw-socket.sh
ExecStartPre=${INSTALL_DIR}/wait-bt-sink.sh
ExecStart=/usr/bin/snapclient -h 127.0.0.1 -p ${SNAPCAST_PORT} --player pulse:server=unix:/run/user/${SOUNDSPOT_UID}/pulse/native
ExecStartPost=/bin/sleep 5
ExecStartPost=${INSTALL_DIR}/play_welcome.sh
Restart=always
RestartSec=5
SyslogIdentifier=soundspot-client
StandardOutput=journal

[Install]
WantedBy=multi-user.target
