# src/templates/soundspot-client-master.service
[Unit]
Description=SoundSpot — Snapclient local (BT speaker)
After=snapserver.service bluetooth.service
Wants=snapserver.service bluetooth.service

[Service]
Type=simple
User=${SOUNDSPOT_USER}
Group=${SOUNDSPOT_USER}
Environment="XDG_RUNTIME_DIR=/run/user/%U"
Environment="PULSE_SERVER=unix:/run/user/%U/pulse/native"
ExecStartPre=/bin/sleep 15
ExecStart=/usr/bin/snapclient -h 127.0.0.1 -p ${SNAPCAST_PORT} --player pulse
ExecStartPost=/bin/sleep 5
ExecStartPost=${INSTALL_DIR}/play_welcome.sh
Restart=always
RestartSec=5
StandardOutput=journal

[Install]
WantedBy=multi-user.target