# src/templates/soundspot-firewall.service
[Unit]
Description=SoundSpot — Pare-feu (iptables + portail captif)
After=ipset-soundspot.service soundspot-ap.service
Requires=ipset-soundspot.service soundspot-ap.service
Before=hostapd.service dnsmasq.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/soundspot/soundspot-firewall.sh
ExecStop=/sbin/iptables -F FORWARD
StandardOutput=journal

[Install]
WantedBy=multi-user.target
