# src/templates/soundspot-client-satellite.service
[Unit]
Description=SoundSpot Satellite — Snapclient
After=bluetooth.service network-online.target time-sync.target user@${SOUNDSPOT_UID}.service
Wants=bluetooth.service network-online.target

[Service]
Type=simple
EnvironmentFile=${INSTALL_DIR}/soundspot.conf
User=${SOUNDSPOT_USER}
Environment="XDG_RUNTIME_DIR=/run/user/${SOUNDSPOT_UID}"
Environment="PULSE_SERVER=unix:/run/user/${SOUNDSPOT_UID}/pulse/native"
ExecStartPre=${INSTALL_DIR}/wait-pw-socket.sh
ExecStartPre=${INSTALL_DIR}/wait-bt-sink.sh
ExecStart=/usr/bin/snapclient -h ${MASTER_HOST} -p ${SNAPCAST_PORT} --player pulse:buffer_time=200
Restart=always
RestartSec=10
StandardOutput=journal

[Install]
WantedBy=multi-user.target
