# src/config/services/soundspot-battery.service
[Unit]
Description=SoundSpot — Monitoring batterie solaire (INA219)
After=soundspot-presence.service

[Service]
Type=simple
EnvironmentFile=${INSTALL_DIR}/soundspot.conf
Environment="XDG_RUNTIME_DIR=/run/user/${SOUNDSPOT_UID}"
# On pointe directement vers le venv ~/.astro
ExecStart=/home/${SOUNDSPOT_USER}/.astro/bin/python3 ${INSTALL_DIR}/backend/system/battery_monitor.py
Restart=on-failure
RestartSec=60
User=${SOUNDSPOT_USER}
SyslogIdentifier=soundspot-battery
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target