# src/config/services/soundspot-decoder.service
[Unit]
Description=SoundSpot — Décodeur Icecast → Snapcast (ffmpeg)
After=icecast2.service snapserver.service
Wants=icecast2.service snapserver.service

[Service]
Type=simple
ExecStartPre=/bin/bash -c 'mkfifo /dev/shm/snapfifo /dev/shm/snapfifo_mic 2>/dev/null; true'
ExecStart=${INSTALL_DIR}/backend/audio/decoder.sh
Restart=always
RestartSec=3
SyslogIdentifier=soundspot-decoder
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
