# src/config/services/soundspot-fleet-relay.service
[Unit]
Description=SoundSpot — Relay NOSTR flotte (port 29999, cluster-local)
After=network.target snapserver.service

[Service]
Type=simple
EnvironmentFile=${INSTALL_DIR}/soundspot.conf
Environment="FLEET_RELAY_PORT=29999"
ExecStart=/usr/bin/python3 ${INSTALL_DIR}/backend/system/fleet_relay.py
Restart=always
RestartSec=5
SyslogIdentifier=fleet-relay
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
